# targets/fs.src/etc/inittab
::sysinit:/etc/init.d/rcS
#::respawn:-/bin/sh
::askfirst:-/bin/login
#::ctrlaltdel:/bin/umount -a -r

